# systemd/hanabi-gpu-server@.service
[Unit]
Description=Hanabi GPU Training Server
After=network.target
Documentation=https://github.com/yourusername/hanabi

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/%i/hanabi
Environment="PATH=/home/%i/hanabi/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/default/hanabi-gpu-server

# Main command
ExecStart=/home/%i/hanabi/.venv/bin/python -m rl_hanabi.training.distributed.gpu_server \
    --host 0.0.0.0 \
    --port ${GPU_PORT:-5555} \
    --auth-token ${AUTH_TOKEN} \
    --device auto \
    --checkpoint-dir /home/%i/hanabi/checkpoints

# Restart settings
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Resource limits
MemoryMax=16G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hanabi-gpu

[Install]
WantedBy=multi-user.target
